# bluelog.crontab
*/2 * * * * root /etc/init.d/bluelog status
